# infra/hetzner/templates/suno-assistant-<job>.timer.j2
[Unit]
Description=suno-assistant <job> schedule

[Timer]
OnCalendar=hourly
Persistent=true
Unit=suno-assistant-<job>.service

[Install]
WantedBy=timers.target
